# systemd/timelapse-compile-daily.service
[Unit]
Description=Compile yesterday's birdsnest timelapse

[Service]
Type=oneshot
ExecStart=@BINDIR@/timelapse-compile-daily
